# retroroot/board/x86_64/kms/overlay/usr/lib/systemd/system/retroarch.service
[Unit]
Description=retroarch
After=syslog.target network.target
RequiresMountsFor=/data

[Service]
ExecStart=/bin/retroarch -c /data/retroarch/config/retroarch.cfg
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
